# uptime-kuma/rootfs/etc/s6-overlay/s6-rc.d/discovery/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Uptime Kuma
# Sends discovery information to Home Assistant.
# ==============================================================================
declare config

# Wait for Uptime Kuma to become available
bashio::net.wait_for 3001 127.0.0.1 300

config=$(bashio::var.json \
    url "http://$(hostname):3001"
)

if bashio::discovery "uptime_kuma" "${config}" > /dev/null; then
    bashio::log.info "Successfully sent discovery information to Home Assistant."
else
    bashio::log.error "Failed to send discovery information to Home Assistant!"
fi
