# uptime-kuma/rootfs/etc/s6-overlay/s6-rc.d/uptime-kuma/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Uptime Kuma
# Starts Uptime Kuma
# ==============================================================================
export NODE_PATH=/opt/node_modules

cd /opt/uptime-kuma || bashio::exit.nok "Could not change directory to Uptime Kuma"

bashio::log.info "Starting Uptime Kuma..."

exec node server/server.js --data-dir="/data/"
